# components/tau_v1_1/agentes/rho/rho_agent.service
[Unit]
Description=Agente de Percepção TAU (RHO)
After=network.target

[Service]
ExecStart=/usr/local/bin/rho_agent
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=rho_agent

[Install]
WantedBy=multi-user.target
