# components/phase_redistributor/tzinor-gateway.service
[Unit]
Description=Tzinor Gateway Phase Redistributor
After=network.target

[Service]
Type=simple
User=arkhe
WorkingDirectory=/usr/local/bin
ExecStart=/usr/local/bin/tzinor_gateway --config /etc/tzinor/gateway_config.yaml
Restart=on-failure
RestartSec=5

# Resources
CPUWeight=100
MemoryLimit=1G

[Install]
WantedBy=multi-user.target
